# distro/dmapd.service
[Unit]
Description=DMAP Service

[Service]
ExecStart=/usr/sbin/dmapd -f

[Install]
WantedBy=multi-user.target
